# ansible-playbook/roles/clamav/files/ansible.cron
0 3 * * * root /usr/bin/freshclam ; /usr/bin/clamscan / --recursive=yes -i | mail -s clamav_log_BoostedNami support@boostedshop.com -